# dev-db/elasticsearch/files/elasticsearch-initd
#!/sbin/runscript

ESEARCH_PID=/var/run/elasticsearch/elasticsearch.pid

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --quiet \
		--user elasticsearch:elasticsearch \
		--exec /opt/bin/elasticsearch -- \
			-p ${ESEARCH_PID}
	
	eend $? || return 1
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop \
		--pidfile ${ESEARCH_PID}
	
	eend $?
}
